# bin/linux/hdlj.web.captcha.test1.service
[Unit]
Description=hdlj.web.captcha.test1.service
After=network.target redis consul
Wants=redis consul
Conflicts=redis consul

[Service]
User=roisoft
Group=roisoft

Type=simple
#EnvironmentFile=/usr/local/bin/hdlj/web_server.cfg
ExecStart=/usr/local/bin/hdlj/web_server
ExecReload=/bin/kill -TERM $MAINPID
KillMode=process
Restart=on-failure
RestartSec=5
PIDFile=/var/run/hdlj.web.captcha.test1.pid
#PrivateTmp=true

LimitCORE=infinity
LimitNOFILE=102400
LimitNPROC=102400

[Install]
WantedBy=multi-user.target
