# bin/linux/hdlj.game.test1.service
[Unit]
Description=hdlj.game.test1.service
After=network.target redis consul hdlj.db.0 hdlj.master.test1
Wants=redis consul hdlj.db.0 hdlj.master.test1
Conflicts=redis consul

[Service]
User=roisoft
Group=roisoft

Type=simple
#EnvironmentFile=/usr/local/bin/hdlj/game_server.cfg
ExecStart=/usr/local/bin/hdlj/game_server
ExecReload=/bin/kill -TERM $MAINPID
KillMode=process
Restart=on-failure
RestartSec=5
PIDFile=/var/run/hdlj.game.test1.pid
#PrivateTmp=true

LimitCORE=infinity
LimitNOFILE=102400
LimitNPROC=102400

[Install]
WantedBy=multi-user.target